# openclaw/launchagents/ai.openclaw.8sleep-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.8sleep-snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dbochman/.openclaw/bin/8sleep-snapshot.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dbochman/.openclaw/bin:/opt/homebrew/bin:/usr/local/bin:/usr/sbin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/dbochman</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>50</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/dbochman/.openclaw/logs/8sleep-snapshot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dbochman/.openclaw/logs/8sleep-snapshot.err.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/dbochman</string>
</dict>
</plist>
