# openclaw/launchagents/ai.openclaw.usage-token-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.usage-token-push</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>security find-generic-password -s "Claude Code-credentials" -w 2>/dev/null | ssh -o IdentityAgent=none -i ~/.ssh/id_mini_to_mbp dylans-mac-mini 'cat > ~/.openclaw/.anthropic-oauth-cache &amp;&amp; chmod 600 ~/.openclaw/.anthropic-oauth-cache'</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dbochman/.openclaw/bin:/opt/homebrew/bin:/usr/local/bin:/usr/sbin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/dbochman</string>
    </dict>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/dbochman/.openclaw/logs/usage-token-push.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dbochman/.openclaw/logs/usage-token-push.err.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/dbochman</string>
</dict>
</plist>
