# openclaw/launchagents/com.openclaw.presence-receive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.presence-receive</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dbochman/.openclaw/workspace/scripts/presence-receive.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dbochman/.openclaw/bin:/opt/homebrew/bin:/usr/local/bin:/usr/sbin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/dbochman</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/dbochman/.openclaw/logs/presence-receive.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dbochman/.openclaw/logs/presence-receive.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
