# openclaw/com.openclaw.bt-connect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.bt-connect</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/ln</string>
        <string>-sf</string>
        <string>/Users/dbochman/.openclaw/bin/bt_op.command</string>
        <string>/tmp/bt_op.command</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
